# launchd/com.nightwatch.transcribe.plist
<?xml version="1.0" encoding="UTF-8"?>
<!-- Transcribe + alert (KeepAlive). Installed by ./install.sh. -->
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key><string>com.nightwatch.transcribe</string>
  <key>ProgramArguments</key>
  <array>
    <string>/Users/alexhedtke/Documents/night-watch/.venv/bin/python</string>
    <string>/Users/alexhedtke/Documents/night-watch/scanner_transcribe.py</string>
  </array>
  <key>WorkingDirectory</key><string>/Users/alexhedtke/Documents/night-watch</string>
  <key>KeepAlive</key><true/>
  <key>StandardOutPath</key><string>/Users/alexhedtke/Documents/night-watch/logs/transcribe.out.log</string>
  <key>StandardErrorPath</key><string>/Users/alexhedtke/Documents/night-watch/logs/transcribe.err.log</string>
</dict>
</plist>
